# scripts/myssh.conf
[program:myssh]
command=/usr/sbin/sshd -D
autostart=true
autorestart=true
user=root
stderr_logfile=/var/log/ssh.err.log
stdout_logfile=/var/log/ssh.out.log